# deploy/cup-public-https.service
[Unit]
Description=CUP direct public HTTPS proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/softcup/to_public/caddy run --config /home/softcup/to_public/Caddyfile --adapter caddyfile
ExecReload=/home/softcup/to_public/caddy reload --config /home/softcup/to_public/Caddyfile --adapter caddyfile
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
