# deploy/cup-huawei-tunnel.service
[Unit]
Description=CUP Huawei reverse SSH tunnel
After=network-online.target cup-api-stack.service
Wants=network-online.target cup-api-stack.service

[Service]
Type=simple
ExecStart=/usr/bin/bash /home/gmn/codes/cup/deploy/watch_huawei_tunnel.sh
Restart=always
RestartSec=3
TimeoutStopSec=20
KillMode=control-group

[Install]
WantedBy=default.target
